# src2/db2fnm
#!/bin/sh
#
# Start/stop the db2fnm daemon
#
### BEGIN INIT INFO
# Provides:          db2fnm
# Required-Start:    $remote_fs $syslog $time postgresql
# Required-Stop:     $remote_fs $syslog $time postgresql
# Should-Start:      $network $named
# Should-Stop:       $network $named
# Default-Start:     2 3 4 5
# Default-Stop:
# Short-Description: Regular background program processing daemon
# Description:       The database to ExaBGP daemon
### END INIT INFO

PATH=/bin:/usr/bin:/sbin:/usr/sbin:/opt/db2dps
DESC="db2fnm daemon"
NAME=db2fnm
PREFIX=/opt/db2dps
DAEMON=${PREFIX}/bin/fnmcfg
SCRIPTNAME=/etc/init.d/"$NAME"

DBINI=${PREFIX}/etc/fnmcfg.ini

test -f $DAEMON || exit 0

. /lib/lsb/init-functions

[ -r ${PREFIX}/etc/default/db2fnm ] && . ${PREFIX}/etc/default/db2fnm

# Read the system's locale and set db2fnm's locale. This is only used for
# setting the charset of mails generated by db2fnm. To provide locale
# information to tasks running under db2fnm, see /etc/pam.d/db2fnm.
#
# We read /etc/environment, but warn about locale information in
# there because it should be in /etc/default/locale.
parse_environment () 
{
    for ENV_FILE in /etc/environment /etc/default/locale; do
        [ -r "$ENV_FILE" ] || continue
        [ -s "$ENV_FILE" ] || continue

         for var in LANG LANGUAGE LC_ALL LC_CTYPE; do
             value=`egrep "^${var}=" "$ENV_FILE" | tail -n1 | cut -d= -f2`
             [ -n "$value" ] && eval export $var=$value

             if [ -n "$value" ] && [ "$ENV_FILE" = /etc/environment ]; then
                 log_warning_msg "/etc/environment has been deprecated for locale information; use /etc/default/locale for $var=$value instead"
             fi
         done
     done

# Get the timezone set.
    if [ -z "$TZ" -a -e /etc/timezone ]; then
        TZ=`cat /etc/timezone` 
    fi
}

get_pid()
{
	BIN=`basename $DAEMON`
	PID=`ps -C $BIN -o pid=`
}

# Parse the system's environment
#if [ "$READ_ENV" = "yes" ] ; then
#    parse_environment
#fi

STOPFILE=`sed '/^ENDLOOP/!d; s/.*=//' ${DAEMON}`
MAX_TRY_SEC=20
EXTRA_OPTS="-l -v"

echo=echo
case ${N}$C in
    "") if $echo "\c" | grep c >/dev/null 2>&1; then
		N='-n'
	else
		C='\c'
	fi ;;
esac

case "$1" in
    start)	log_daemon_msg "Starting db2fnm export db to fnm" ""
            get_pid
            if [ "$PID" != "" ]; then
                echo 
                echo "$BIN allready running on pid $PID"
                RETVAL=1
            else
                start_daemon $DAEMON $EXTRA_OPTS
            fi
            log_end_msg $?
            ;;
    stop)	log_daemon_msg "Stopping db2fnm export db to fnm" ""
            touch ${STOPFILE}
            for s in 0 `seq ${MAX_TRY_SEC}`
            do
                echo $N "$s $C"
                get_pid
                if [ "$PID" != "" ]; then
                    sleep 1
                else
                    break
                fi
            done
            if [ "$PID" != "" ]; then
                echo "$N killing $PID $C"
                kill -9 $PID
                RETVAL=$?
            else
                RETVAL=0
            fi
            log_end_msg $RETVAL
            /bin/rm -f ${STOPFILE}
            ;;
    reload|force-reload|restart) log_daemon_msg "Restarting DDoS detector" "db2dps" 
            $0 stop
            $0 start
            ;;
    status)
            status_of_proc $DAEMON $NAME && exit 0 || exit $?
            ;;
    *)	log_action_msg "Usage: /etc/init.d/db2fnm {start|stop|status|restart|reload|force-reload}"
            exit 2
            ;;
esac

exit 0

