# src/FreeBSD-i2dps_1.0-14/etc/cron.d/watchdog
*/5 * * * * root [ -x /opt/i2dps/bin/watchdog.sh ] && /opt/i2dps/bin/watchdog.sh
